# jhandubamb/etc/systemd/system/snap-slack-79.mount
[Unit]
Description=Mount unit for slack, revision 79
After=snapd.mounts-pre.target
Before=snapd.mounts.target
Before=local-fs.target

[Mount]
What=/var/lib/snapd/snaps/slack_79.snap
Where=/snap/slack/79
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
