# jhandubamb/etc/systemd/system/snap-postman-206.mount
[Unit]
Description=Mount unit for postman, revision 206
After=snapd.mounts-pre.target
Before=snapd.mounts.target
Before=local-fs.target

[Mount]
What=/var/lib/snapd/snaps/postman_206.snap
Where=/snap/postman/206
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
