# jhandubamb/etc/systemd/system/snap-code-132.mount
[Unit]
Description=Mount unit for code, revision 132
After=snapd.mounts-pre.target
Before=snapd.mounts.target
Before=local-fs.target

[Mount]
What=/var/lib/snapd/snaps/code_132.snap
Where=/snap/code/132
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
